# kcptun.conf
[program:kcptun]
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=~/kcptun.log