# kcptun.conf
[program:kcptun]
priority=50
startsecs=0
startretries=0
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/kcptun.log
